# mail.kmx.io/etc/rc.d/mail_kmx_io
#!/bin/ksh

daemon="/bin/sh /home/exwm_web/start_"
daemon_user=exwm_web
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
