# skills.kmx.io/etc/rc.d/skills_kmx_io
#!/bin/ksh

daemon="/bin/sh /home/skills/start_"
daemon_user=skills
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
