# git.kmx.io/etc/rc.d/git_test_kmx_io
#!/bin/ksh

daemon="/bin/sh /home/git_test/_/start_"
daemon_user=git_test
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
