# skills/etc/rc.d/skills_staging
#!/bin/ksh

daemon="/bin/sh /home/skills-staging/start_"
daemon_user=skills-staging
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
