# metrics.kmx.io/etc/rc.d/metrics_kmx_io
#!/bin/ksh

daemon="/bin/sh /home/kmxmetrics/start_"
daemon_user=kmxmetrics
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
