# metrics.kmx.io/etc/rc.d/metrics_kmx_io
#!/bin/ksh

daemon="/bin/sh /home/metrics_kmx_io/start_"
daemon_user=metrics_kmx_io
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
