# www.kmx.io/OpenBSD/etc/rc.d/www_kmx_io
#!/bin/ksh

daemon="/bin/sh /home/www.kmx.io/start_"
daemon_user=www_kmx_io
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
