# git.kmx.io/OpenBSD/etc/rc.d/fcgiwrap
#!/bin/ksh

daemon="/usr/local/sbin/fcgiwrap -c 12 -s tcp:127.0.0.1:9001"
daemon_user=git
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
