# git.kmx.io/OpenBSD/etc/rc.d/git_kmx_io
#!/bin/ksh

daemon="/bin/sh /home/git/_/start_"
daemon_user=git

. /etc/rc.d/rc.subr

pexp=".*/beam[.]smp .* -home /home/git .*"

rc_cmd $1
