# exwm/etc/rc.d/exwm_smtp
#!/bin/ksh

daemon="/bin/sh /home/exwm_smtp/start_"
daemon_user=exwm_smtp
rc_bg=YES

. /etc/rc.d/rc.subr

rc_cmd $1
